# systemd/whisper-en.service
[Unit]
Description=Whisper EN STT server (Vulkan iGPU, whisper.cpp-amd)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rahlquist
WorkingDirectory=/home/rahlquist/whisper-models/whisper.cpp-amd
ExecStart=/home/rahlquist/whisper-models/whisper.cpp-amd/build/bin/whisper-server --model /home/rahlquist/whisper-models/kb-whisper-large-ggml/ggml-model-medium.en.bin --port 8083 --convert
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
